# skills/levnikolaevich--claude-code-skills/skills-catalog/ln-030-vps-bootstrap/references/agent-update.service
[Unit]
Description=Nightly system-wide agent toolchain update (Claude/Codex CLI, marketplace, plugins) + restart of all god-services
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/agent-update
TimeoutStartSec=30min
StandardOutput=journal
StandardError=journal
