# skills/levnikolaevich--claude-code-skills/skills-catalog/ln-030-vps-bootstrap/references/claude-relay-bot.service
[Unit]
Description=Claude relay bot - Telegram bridge + per-user sessions UI for ${SERVICE_PREFIX}
After=network-online.target
Wants=network-online.target
# Relay starts ${SERVICE_PREFIX}-god@<telegram_user_id>.service on demand and
# keeps queued inbound Telegram text durable until that user's tmux pane is ready.

[Service]
Type=simple
User=${BOT_USER}
Group=${BOT_USER}
EnvironmentFile=/etc/${PROJECT_NAME}/secrets.env
Environment=PROJECT_NAME=${PROJECT_NAME}
Environment=PROJECT_DIR=${PROJECT_DIR}
Environment=SERVICE_PREFIX=${SERVICE_PREFIX}
Environment=BOT_USER=${BOT_USER}
Environment=RELAY_HOOK_PORT=${RELAY_HOOK_PORT}
WorkingDirectory=/opt/${SERVICE_PREFIX}-relay-bot
ExecStart=/bin/bash -lc '. /home/${BOT_USER}/.nvm/nvm.sh && exec node dist/index.js'
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal

MemoryMax=256M
CPUQuota=50%

# Required for sudo -n systemctl start/is-active ${SERVICE_PREFIX}-god@*.service.
# Bootstrap must install a sudoers fragment restricted to this service prefix.
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/${PROJECT_NAME} /tmp ${PROJECT_DIR}/.agent-home/users
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallArchitectures=native
LockPersonality=true
RestrictRealtime=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
# Do not enable PrivateTmp: relay shares /tmp/tmux-<UID> with god@* panes.

[Install]
WantedBy=multi-user.target
