# skills/levnikolaevich--claude-code-skills/skills-catalog/ln-030-vps-bootstrap/references/claude-relay-bot.service
[Unit]
Description=Claude relay bot — Telegram bridge + sessions UI for ${SERVICE_PREFIX} god-session
After=network-online.target ${SERVICE_PREFIX}-god.service
Wants=network-online.target
# Relay is intentionally independent from ${SERVICE_PREFIX}-god.service. It must
# survive god-session restarts triggered by /new_session or [Resume] and keep
# queued inbound Telegram text durable until tmux is ready again.

[Service]
Type=simple
User=${BOT_USER}
Group=${BOT_USER}
EnvironmentFile=/etc/${PROJECT_NAME}/secrets.env
ExecStart=/home/${BOT_USER}/.venv-relay/bin/python /usr/local/bin/claude-relay-bot.py
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal

# Resource caps; relay is a thin polling loop with light fs scans.
MemoryMax=128M
CPUQuota=50%

# Hardening: keep the Telegram ingress process tightly sandboxed.
# Mitigates blast radius if the bot token leaks: the relay process is the
# only path from Telegram into the god-session, so it must be tightly sandboxed.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
# State dir is the single writable location (atomic command file, sessions-dir.path,
# last-god-error.json, relay.db). /tmp/tmux-<UID> needed for tmux socket.
ReadWritePaths=/var/lib/${PROJECT_NAME} /tmp
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallArchitectures=native
LockPersonality=true
RestrictRealtime=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
# IMPORTANT: do NOT enable PrivateTmp — relay needs to share /tmp/tmux-<UID>
# socket with the ${SERVICE_PREFIX}-god service to send-keys into the same pane.

[Install]
WantedBy=multi-user.target
