# skills/levnikolaevich--claude-code-skills/skills-catalog/ln-030-vps-bootstrap/references/dispatch.timer
[Unit]
Description=Hourly /${DISPATCH_COMMAND_NAME} trigger for ${SERVICE_PREFIX} god-session (durable, replaces in-session /loop)

[Timer]
OnCalendar=*-*-* *:07:00
Persistent=true
Unit=${SERVICE_PREFIX}-dispatch.service

[Install]
WantedBy=timers.target
