# skills/levnikolaevich--claude-code-skills/skills-catalog/ln-030-vps-bootstrap/references/dispatch.service
[Unit]
Description=Inject /${DISPATCH_COMMAND_NAME} into ${SERVICE_PREFIX}-god tmux session (replaces internal /loop)
Requires=${SERVICE_PREFIX}-god.service
After=${SERVICE_PREFIX}-god.service

[Service]
Type=oneshot
User=${BOT_USER}
Group=${BOT_USER}
ExecStart=/usr/bin/tmux send-keys -t ${SERVICE_PREFIX}-god "/${DISPATCH_COMMAND_NAME}" Enter
StandardOutput=journal
StandardError=journal
