# skills/levnikolaevich--claude-code-skills/plugins/setup-environment/skills/ln-030-vps-bootstrap/references/templates/agent-update.timer
[Unit]
Description=Nightly system-wide agent toolchain update timer (one-per-VPS, fires for all projects)

[Timer]
OnCalendar=*-*-* 03:37:00
RandomizedDelaySec=20m
Persistent=true

[Install]
WantedBy=timers.target
