# skills/levnikolaevich--claude-code-skills/plugins/setup-environment/skills/ln-030-vps-bootstrap/references/templates/dispatch.timer
[Unit]
Description=15-minute ${SERVICE_PREFIX} hex-relay task polling trigger

[Timer]
OnCalendar=*-*-* *:0/15:00
Persistent=true
Unit=${SERVICE_PREFIX}-dispatch.service

[Install]
WantedBy=timers.target
