# skills/levnikolaevich--claude-code-skills/plugins/setup-environment/skills/ln-030-vps-bootstrap/references/templates/dispatch.service
[Unit]
Description=Poll provider tasks through ${SERVICE_PREFIX} hex-relay control plane
Requires=${SERVICE_PREFIX}-hex-relay.service
After=${SERVICE_PREFIX}-hex-relay.service

[Service]
Type=oneshot
User=${BOT_USER}
Group=${BOT_USER}
EnvironmentFile=/etc/${PROJECT_NAME}/secrets.env
ExecStart=/usr/bin/curl -fsS -X POST -H "Authorization: Bearer ${RELAY_HTTP_TOKEN}" http://127.0.0.1:${RELAY_HOOK_PORT}/tasks/poll
StandardOutput=journal
StandardError=journal
