# skills/levnikolaevich--claude-code-skills/plugins/setup-environment/skills/ln-030-vps-bootstrap/references/templates/god-session-codex.service
[Unit]
Description=${PROJECT_NAME} god-session-codex for Telegram user %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=${BOT_USER}
Group=${BOT_USER}
Environment=HOME=/home/${BOT_USER}
Environment=PROJECT_NAME=${PROJECT_NAME}
Environment=PROJECT_DIR=${PROJECT_DIR}
Environment=SERVICE_PREFIX=${SERVICE_PREFIX}
Environment=BOT_USER=${BOT_USER}
Environment=AGENT_SKILLS_DIR=${AGENT_SKILLS_DIR}
Environment=AGENT_HOME=${PROJECT_DIR}/.agent-home/users/%i
Environment=OPERATOR_USER_ID=%i
Environment=RELAY_HOOK_PORT=${RELAY_HOOK_PORT}
WorkingDirectory=${PROJECT_DIR}
ExecStart=/usr/local/bin/${SERVICE_PREFIX}-god-codex %i
ExecStop=/usr/bin/tmux -L ${SERVICE_PREFIX} kill-session -t ${SERVICE_PREFIX}-god-codex-%i
Restart=always
RestartSec=15
TimeoutStopSec=30
StandardOutput=append:/var/log/${PROJECT_NAME}-god.log
StandardError=append:/var/log/${PROJECT_NAME}-god.log

# Project state contains per-user codex resume state under users/%i/.
StateDirectory=${PROJECT_NAME}
StateDirectoryMode=0700

MemoryMax=2G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
